# deploy/smart-home-system.service
[Unit]
Description=Smart Home System
Documentation=https://github.com/stanleylabs/smart-home-system
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=smarthub
Group=smarthub
WorkingDirectory=/opt/smart-home-system
ExecStartPre=/bin/sh -c 'iptables -t nat -C PREROUTING -p tcp --dport 80 -j REDIRECT --to-port 3000 2>/dev/null || iptables -t nat -A PREROUTING -p tcp --dport 80 -j REDIRECT --to-port 3000'
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5

# Graceful shutdown via SIGTERM (handled by the hub)
KillSignal=SIGTERM
TimeoutStopSec=15

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/smart-home-system/data /opt/smart-home-system/config

# CAP_NET_RAW = BLE HCI sockets for Matter BLE commissioning
AmbientCapabilities=CAP_NET_RAW

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
